# content/etc/init.d/cloudflared
#!/sbin/openrc-run
name="cloudflared"
description="Cloudflare Tunnel"

etc_dir="/etc/$name"
lib_dir="/var/lib/$name"
log_dir="/var/log/$name"
run_dir="/run/$name"

username="cloudflare"

pidfile="$run_dir/$name.pid"
logfile="$log_dir/$name.log"

command="/usr/bin/$name"
command_user="$username:$username"

depend() {
    need net
}

start_pre() {
    # Ensure the user exists
    if ! id -u $username >/dev/null 2>&1; then
        echo "Creating $username user..."
        addgroup $username
        adduser -S -D -H -G $username $username
    fi

    checkpath -d -o $command_user -m755 \
            $etc_dir

    checkpath -f -o $command_user -m644 \
            $config \
            $logfile
}

start() {
    ebegin "Starting $name"
    start-stop-daemon --start --quiet --pidfile $pidfile --exec $command -- $command_args
    eend $?
}

stop() {
    ebegin "Stopping $name"
    start-stop-daemon --stop --quiet --pidfile $pidfile
    eend $?
}
